# 90system-upgrade/system-upgrade.target
[Unit]
Description=System Upgrade
Documentation=man:system-upgrade.target(7)
Requires=sysinit.target sockets.target
After=sysinit.target sockets.target
AllowIsolate=yes
